# deploy/knwn4-agent.service
[Unit]
Description=KNWN4 Content Agent
After=network.target

[Service]
Type=simple
User=jah
WorkingDirectory=/opt/knwn4-agent
ExecStart=/usr/bin/doppler run --project knwn4-claude-agent --config prd -- /opt/knwn4-agent/.venv/bin/python3 -m src.main
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
